# launchctl/Library/LaunchAgents/com.asmvik.yabai.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Generated by the yabai "install-service" flag, then hand-edited: PATH is
  trimmed to a stable, machine-agnostic set (the generator snapshots the
  invoking shell's PATH, baking in session-specific junk). Manage the service
  with the "start-service"/"stop-service"/"restart-service" flags only; never
  re-run "install-service", as through the ~/Library/LaunchAgents stow symlink
  it would overwrite this tracked file.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.asmvik.yabai</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/yabai</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
 	     <false/>
 	     <key>Crashed</key>
 	     <true/>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/yabai_mvshmakov.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/yabai_mvshmakov.err.log</string>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>Nice</key>
    <integer>-20</integer>
</dict>
</plist>