# _archive/launchctl/Library/LaunchAgents/com.mvshmakov.setenv.GNUPGHOME.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <!--
            We need to "launchctl setenv" mostly because this is the only way to make GPG Tools
            work correctly the XDG Base Directory spec-compliant GNUPGHOME.
            Source: https://gpgtools.tenderapp.com/discussions/problems/49849-gpg-keychain-non-standard-location-for-gnupg-gnupghome-directory
        -->
        <key>Label</key>
        <string>com.mvshmakov.setenv.GNUPGHOME</string>
        <key>ProgramArguments</key>
        <array>
            <string>/bin/launchctl</string>
            <string>setenv</string>
            <string>GNUPGHOME</string>
            <string>/Users/mvshmakov/.local/share/gnupg</string>
        </array>
        <key>RunAtLoad</key>
        <true />
    </dict>
</plist>